# distrib/headless/debinstall/deb-src/DEBIAN/postinst
#!/bin/bash

update-rc.d maracuya-jukebox defaults

service maracuya-jukebox start